# init/gipns.service
[Unit]
Description=Git to IPNS

[Service]
ExecStart=~/go/bin/gipns
Restart=on-failure

[Install]
WantedBy=default.target